# .config/autostart/braveTabFix.desktop
[Desktop Entry]
Type=Application
Exec=brave --no-startup-window
Hidden=false
NoDisplay=false
