# .config/autostart/chromeTabFix.desktop
[Desktop Entry]
Type=Application
Exec=google-chrome-stable --no-startup-window
Hidden=false
NoDisplay=false
