# .config/autostart/imwheel.desktop
[Desktop Entry]
Type=Application
Exec=imwheel 45
Hidden=false
NoDisplay=false
